# systemd/loopia.timer
[Unit]
Description=Timer for Loopia DynDNS

[Timer]
OnUnitActiveSec=1m
OnBootSec=20s

[Install]
WantedBy=timers.target
